# systemd/popub-local@.service.in
[Unit]
Description=Popub Local
After=network-online.target

[Service]
Type=simple
AmbientCapabilities=CAP_NET_BIND_SERVICE
EnvironmentFile=/etc/popub/local.d/%i
EnvironmentFile=%h/.config/popub/%i
ExecStart=<PREFIX>/bin/popub-local "$LOCAL_ADDR" "$RELAY_ADDR" "$AUTH_KEY"
LimitNOFILE=1048576
Restart=always
RestartSec=3
User=<USER>

[Install]
WantedBy=multi-user.target
